# supervisord.conf
[supervisord]
nodeamon=true

[program:apache2]
command = /etc/init.d/apache2 start
user = root
autostart = true
